# automem-sync.service
[Unit]
Description=AutoMem Transcript Sync
Requires=automem.service
After=automem.service

[Service]
Type=simple
WorkingDirectory=/home/clawdbot/clawd/automem
ExecStart=/usr/bin/python3 /home/clawdbot/clawd/automem/scripts/transcript-sync.py --interval 30
Restart=always
RestartSec=10
User=clawdbot
Environment=AUTOMEM_ENDPOINT=http://localhost:8001
Environment=AUTOMEM_API_TOKEN=86d634aaa78b4376bad0935e09a4ee4c76c8cbb30280c60d

[Install]
WantedBy=multi-user.target
